# scripts/systemd/cocktail-hardware.service.template
[Unit]
Description=Cocktail Machine Hardware API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__REPO_ROOT__
ExecStart=/bin/bash __REPO_ROOT__/scripts/start-hardware.sh
Restart=on-failure
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
