# scripts/cocktail-kiosk.service
[Unit]
Description=Cocktail Machine Kiosk Interface
After=cocktail-machine.service graphical-session.target
Wants=cocktail-machine.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/potion-master-pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority

# Start the preview server
ExecStart=/bin/bash -c 'npm run preview -- --host 0.0.0.0 --port 3000'

Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target